# uu_prepare
#!/bin/sh /etc/rc.common

START=99

start() {
  # Configure Network using UCI for safety
  # This avoids syntax errors if variables like UU_LAN_DNS contain spaces
  uci set network.lan.proto='static'
  uci set network.lan.ipaddr="$UU_LAN_IPADDR"
  uci set network.lan.gateway="$UU_LAN_GATEWAY"
  uci set network.lan.netmask="$UU_LAN_NETMASK"
  
  # Handle DNS: Clear old list and add new ones (handles multiple IPs correctly)
  uci -q delete network.lan.dns
  for dns in $UU_LAN_DNS; do
      uci add_list network.lan.dns="$dns"
  done
  
  uci commit network

  # Restart network to apply changes
  /etc/init.d/network restart
  
  # CRITICAL: Wait for network interface to fully come up before restarting dnsmasq
  # This fixes the race condition where dnsmasq starts before the IP is assigned
  sleep 5

  # Configure DHCP and DNS
  uci set dhcp.lan.ignore=1
  uci commit dhcp
  /etc/init.d/dnsmasq restart

  # Ensure UU Plugin is running
  /etc/init.d/uuplugin start
}
